# contrib/session-helper/colord-session.service.in
[Unit]
Description=Color management helper

[Service]
Type=dbus
BusName=org.freedesktop.ColorHelper
ExecStart=@servicedir@/colord-session
