# contrib/colord-sane/colord-sane.service.in
[Unit]
Description=Daemon for monitoring attached scanners and registering them with colord

[Service]
Type=dbus
BusName=org.freedesktop.colord-sane
ExecStart=@servicedir@/colord-sane
User=@daemon_user@
PrivateNetwork=yes
PrivateTmp=yes
