# data/colord.service.in
[Unit]
Description=Manage, Install and Generate Color Profiles

[Service]
Type=dbus
BusName=org.freedesktop.ColorManager
ExecStart=@servicedir@/colord
User=@daemon_user@
# We think that udev's AF_NETLINK messages are being filtered when
# network namespacing is on.
# PrivateNetwork=yes
PrivateTmp=yes
